# bypass-ru/bypass-ru.service
[Unit]
Description=Bypass routing for .ru/.xn--p1ai via main table (mark 0x1)
After=network-pre.target
Before=network.target NetworkManager.service
Wants=network-pre.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/sbin/nft -f /etc/nftables.d/bypass-ru.nft
ExecStart=/usr/sbin/ip    rule add fwmark 0x1 lookup main pref 8000
ExecStart=/usr/sbin/ip -6 rule add fwmark 0x1 lookup main pref 8000
ExecStop=-/usr/sbin/ip    rule del fwmark 0x1 lookup main pref 8000
ExecStop=-/usr/sbin/ip -6 rule del fwmark 0x1 lookup main pref 8000
ExecStop=-/usr/sbin/nft delete table inet bypass

[Install]
WantedBy=multi-user.target
